# archive/edu.tamu.oalservicemanager.login.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>edu.tamu.oalservicemanager.login</string>
    <key>ProgramArguments</key>
    <array>
        <string>python</string>
        <string>-u</string>
        <string>/Library/Scripts/OAL/servicemanager.py</string>
        <string>login</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>EnableTransactions</key>
    <false/>
    <key>ExitTimeOut</key>
    <integer>40</integer>
</dict>
</plist>